# tools/image-targets/pi5-bookworm-armhf/rootfs/etc/systemd/system/PGenerator.service
[Unit]
Description=PGenerator web and pattern services
After=network-online.target rcPGenerator.service
Wants=network-online.target rcPGenerator.service
ConditionPathExists=/etc/init.d/PGenerator

[Service]
Type=oneshot
RemainAfterExit=yes
Environment=LANG=C.UTF-8
Environment=LC_ALL=C.UTF-8
ExecStart=/etc/init.d/PGenerator start
ExecStop=/etc/init.d/PGenerator stop

[Install]
WantedBy=multi-user.target
